# self-docker-image/apacheandsamba/service.conf
[program:apache2]
command=/usr/sbin/apache2ctl -D FOREGROUND

[program:samba]
command=/usr/sbin/service samba start
